# tuned-config/debian/bullseye/deb_folder/postinst
#!/bin/sh

set -e

systemctl enable tuned.service > /dev/null 2>&1 ||:

#DEBHELPER#
