# nfs-utils-config/debian/all/deb_folder/postinst
#!/bin/sh

set -e

systemctl unmask nfs-common.service
systemctl enable nfs-common.service > /dev/null 2>&1 || :

#DEBHELPER#
