# systemd-config/source-debian/utils.slice
[Unit]
Description=StarlingX utility services slice
Before=slices.target

[Slice]
MemoryAccounting=true
CPUAccounting=true
CPUShares=128
